# provision/ansible/roles/prereq/files/containerd-cleanup.timer
[Unit]
Description=containerd cleanup timer

[Timer]
OnUnitInactiveSec=12h

[Install]
WantedBy=timers.target
